# mac/com.vihpol.scandemo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.vihpol.scandemo</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/node</string>
      <string>server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vihpol/Documents/Codex/2026-05-20/im-in-a-networks-company-and</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>PORT</key>
      <string>5173</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/vihpol/Documents/Codex/2026-05-20/im-in-a-networks-company-and/scan-demo.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vihpol/Documents/Codex/2026-05-20/im-in-a-networks-company-and/scan-demo.err.log</string>
  </dict>
</plist>
